# gensyn.service
#!/bin/bash

LOGNAME=$(logname 2>/dev/null)
USER_HOME="/home/$LOGNAME"

# Tao service rl-swarm
echo "Tao service rl-swarm..."
cat << 'EOF' > /etc/systemd/system/rl-swarm.service
[Unit]
Description=RL Swarm Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/rl-swarm
ExecStart=/bin/bash -c 'source /root/rl-swarm/.venv/bin/activate && ./run_rl_swarm.sh'
Restart=always
RestartSec=5
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
EOF

# Tao service rl-swarm2
echo "Tao service rl-swarm2..."
cat << EOF > /etc/systemd/system/rl-swarm2.service
[Unit]
Description=Auto Run Gensyn Script
After=network.target

[Service]
User=root
ExecStart=$USER_HOME/gensyn/run/auto.run
Restart=always
RestartSec=10
Environment=LOGNAME=$LOGNAME

[Install]
WantedBy=multi-user.target
EOF

# Tao thu muc neu chua co
echo "Tao thu muc $USER_HOME/gensyn/run..."
mkdir -p "$USER_HOME/gensyn/run"

# Tao file auto.run (ghi de neu da ton tai)
echo "Tao file auto.run..."
cat << EOF > "$USER_HOME/gensyn/run/auto.run"
#!/bin/bash

# Thu muc chua du lieu goc
BASE_DIR="$USER_HOME/gensyn/run"

# Thu muc dich cua rl-swarm
TARGET_DIR="/root/rl-swarm"

# Thu muc con de copy temp-data
LOGIN_SUBDIR="modal-login"

# Ham copy file va restart
copy_and_restart() {
    local idx="\$1"
    PEM_SOURCE="\$BASE_DIR/\$idx/swarm.pem"
    TEMP_SOURCE="\$BASE_DIR/\$idx/temp-data"
    
    systemctl stop rl-swarm
    sleep 5
    
    pkill -f "yarn start" 2>/dev/null || true
    pkill -f "node.*modal-login" 2>/dev/null || true
    sleep 3
    
    rm -f "\$TARGET_DIR/swarm.pem"
    cp -f "\$PEM_SOURCE" "\$TARGET_DIR/"
    
    rm -rf "\$TARGET_DIR/\$LOGIN_SUBDIR/temp-data"
    cp -r "\$TEMP_SOURCE" "\$TARGET_DIR/\$LOGIN_SUBDIR/"
    
    systemctl start rl-swarm
    
    # Cho userData.json xuat hien
    while [ ! -f "\$TARGET_DIR/\$LOGIN_SUBDIR/temp-data/userData.json" ]; do
        sleep 5
    done
    
    # Lay ORG_ID
    ORG_ID=\$(awk 'BEGIN { FS = "\"" } !/^[ \t]*[{}]/ { print \$(NF - 1); exit }' \\
        "\$TARGET_DIR/\$LOGIN_SUBDIR/temp-data/userData.json")
    
    # Cho API key kich hoat
    while true; do
        STATUS=\$(curl -s "http://localhost:3000/api/get-api-key-status?orgId=\$ORG_ID")
        if [[ "\$STATUS" == "activated" ]]; then
            break
        else
            sleep 5
        fi
    done
}

# Lan dau dung folder 1
echo "\$(date) - Khoi dong lan dau, dung folder 1"
copy_and_restart 1

# Bat dau vong lap tu folder 2
CURRENT_INDEX=2

# Thoi gian cho toi da (10 phut)
TIMEOUT_SECONDS=600
last_detect_time=\$(date +%s)

while true; do
    LOG_LAST_30S=\$(journalctl -u rl-swarm --since "30 seconds ago" -o cat)
    current_time=\$(date +%s)
    
    if echo "\$LOG_LAST_30S" | grep -q "Joining round"; then
        last_detect_time=\$current_time
        
        # Reset index neu vuot qua 32 truoc khi dung
        if (( CURRENT_INDEX > 32 )); then
            CURRENT_INDEX=1
        fi
        
        echo "\$(date) - Phat hien Joining round, copy & restart voi folder \$CURRENT_INDEX"
        copy_and_restart "\$CURRENT_INDEX"
        ((CURRENT_INDEX++))
    else
        # Neu qua 10 phut khong phat hien Joining round thi chuyen index tiep
        if (( current_time - last_detect_time >= TIMEOUT_SECONDS )); then
            echo "\$(date) - Khong phat hien Joining round trong 10 phut, chuyen sang thu muc \$CURRENT_INDEX"
            
            # Reset index neu vuot qua 32 truoc khi dung
            if (( CURRENT_INDEX > 32 )); then
                CURRENT_INDEX=1
            fi
            
            copy_and_restart "\$CURRENT_INDEX"
            ((CURRENT_INDEX++))
            last_detect_time=\$current_time
        fi
    fi
    
    sleep 20
done
EOF

# Them quyen thuc thi
chmod +x "$USER_HOME/gensyn/run/auto.run"
echo "Da them quyen thuc thi cho auto.run"

# Reload systemd
systemctl daemon-reload

# Cau hinh crontab cho root
CRON_CMD="0 2,5,8,11,14,17,20,23 * * * systemctl restart rl-swarm2"
CRON_COMMENT="# Auto restart rl-swarm2 every 3 hours"

echo "Cau hinh crontab cho root..."

# Lay crontab hien tai cua root
CURRENT_CRON=$(crontab -l 2>/dev/null)

# Kiem tra xem cron job da ton tai chua
if echo "$CURRENT_CRON" | grep -qF "systemctl restart rl-swarm2"; then
    echo "Cron job da ton tai, xoa va them lai..."
    # Xoa dong cu va them lai
    UPDATED_CRON=$(echo "$CURRENT_CRON" | grep -vF "systemctl restart rl-swarm2" | grep -vF "Auto restart rl-swarm2")
    if [ -z "$UPDATED_CRON" ]; then
        (echo "$CRON_COMMENT"; echo "$CRON_CMD") | crontab -
    else
        (echo "$UPDATED_CRON"; echo ""; echo "$CRON_COMMENT"; echo "$CRON_CMD") | crontab -
    fi
else
    echo "Them cron job moi..."
    if [ -z "$CURRENT_CRON" ]; then
        (echo "$CRON_COMMENT"; echo "$CRON_CMD") | crontab -
    else
        (echo "$CURRENT_CRON"; echo ""; echo "$CRON_COMMENT"; echo "$CRON_CMD") | crontab -
    fi
fi

echo "De kich hoat services, chay cac lenh sau:"
echo "  systemctl enable rl-swarm rl-swarm2 && systemctl restart rl-swarm2"
